# stash/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:stash]
command=/opt/stash/bin/start-stash.sh -fg
